# dotfiles/rgb/rgb.service
[Unit]
Description=OpenRGB control

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb --profile /etc/openrgb/pinkplue.orp
User=nobody
Group=i2c

[Install]
WantedBy=multi-user.target